# packages/tblock/tblock-2.5.1.tar.gz/tblock-2.5.1/assets/services/systemd/tblockd.service
[Unit]
Description=TBlock Daemon
Requires=network-online.target
After=network-online.target 

[Service]
ExecStart=/usr/bin/env tblockd --daemon --config=/etc/tblock.conf --quiet

[Install]
WantedBy=default.target
